# Smolit_AdminBot/rootbot.service
[Unit]
Description=RootBot Autonomous System Administration Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root_adminbot
Group=adminbot_group
WorkingDirectory=/opt/rootbot
ExecStart=/usr/local/bin/rootbot start
ExecStop=/usr/local/bin/rootbot stop
Restart=always
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=30
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target

